# archcraft/.config/autostart/clipit.desktop
[Desktop Entry]
Name=ClipIt
Exec=/usr/bin/clipit
Terminal=false