# archcraft/.config/autostart/flameshot.desktop
[Desktop Entry]
Name=Flameshot
Exec=/usr/bin/flameshot
Terminal=false