# linux/.config/autostart/theme.desktop
[Desktop Entry]
Name=Change to Theme default
Exec=/home/edinei/.config/openbox-themes/themes/default/apply.sh
Terminal=false