# linux/.config/autostart/thunderbird-beta.desktop
[Desktop Entry]
Type=Application
Name=Thunderbird beta
Exec=/usr/bin/thunderbird-beta