# dot_config/autostart/localsend_app.desktop
[Desktop Entry]
Type=Application
Name=localsend_app
Comment=localsend_app startup script
Exec=/opt/localsend/localsend --hidden
StartupNotify=false
Terminal=false
